# systemd/tm4-runner.service
[Unit]
Description=TM4 Runner Service
After=network.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/opt/tm4server
Environment=TM4_BASE_PATH=/var/lib/tm4
Environment=TM4_POLL_INTERVAL_S=3
Environment=TM4_CORE_PATH=/opt/tm4-core
Environment=TM4_AUTONOMY_SCRIPT=/opt/tm4-core/mvp/scripts/run_autonomy_loop.py
Environment=TM4_PYTHON_BIN=/opt/tm4-core/venv/bin/python
ExecStart=/opt/tm4server/venv/bin/python -m tm4server.worker
Restart=always
RestartSec=5
StandardOutput=append:/var/log/tm4/tm4-runner.log
StandardError=append:/var/log/tm4/tm4-runner.log

[Install]
WantedBy=multi-user.target
