# systemd/tm4-runner.service
[Unit]
Description=TM4 Execution Runner
After=network.target

[Service]
Type=simple
User=tm4
Group=tm4
WorkingDirectory=/opt/tm4server
# Environment variables for the worker
Environment=TM4_RUNTIME_ROOT=/var/lib/tm4
Environment=TM4SERVER_REPO_ROOT=/opt/tm4server
Environment=TM4CORE_REPO_ROOT=/opt/tm4-core
# Use the virtualenv python
ExecStart=/opt/tm4server/venv/bin/python -m tm4server.execution.worker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
